# appliance_features/auto_reboot/auto_reboot.cron
# Reboot this appliance every week
0 7 * * SUN /home/pi/bin/reboot-appliance
